# live/eu-west-1/database/ursho.service
[Unit]
Description=ursho service

[Service]
Type=simple
WorkingDirectory=/bin
ExecStart=/bin/ursho

[Install]
WantedBy=multi-user.target